# klipper-mcp.service
# Reference template only — install.sh generates the real unit with the correct
# user and paths (derived from the clone location). Replace the __PLACEHOLDERS__
# if you install the service by hand.
#
#   __USER__       service user (e.g. output of `whoami`)
#   __GROUP__      that user's primary group (e.g. output of `id -gn`)
#   __INSTALL_DIR__ absolute path to the cloned repo (e.g. /home/<user>/klipper-mcp)

[Unit]
Description=Klipper MCP Server
After=network.target moonraker.service
Wants=moonraker.service

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__INSTALL_DIR__
Environment="PATH=__INSTALL_DIR__/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=__INSTALL_DIR__/venv/bin/python server.py
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/klipper-mcp.log
StandardError=append:/var/log/klipper-mcp.log

[Install]
WantedBy=multi-user.target
